# config/systemd/acme-dns.service
[Unit]
Description=acme-dns DNS-01 challenge server (Let's Encrypt wildcard auto-renewal)
Documentation=https://github.com/joohoi/acme-dns
Documentation=file:///root/mindGarden/docs/project-management/2026-05-28/SSL_WILDCARD_ACMEDNS_AUTO_RENEW_HANDOFF.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=acme-dns
Group=acme-dns
ExecStart=/usr/local/bin/acme-dns -c /etc/acme-dns/config.cfg
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30
LimitNOFILE=4096

# 53 포트 바인딩 (비특권 계정에서)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# 파일 시스템 보호
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/var/lib/acme-dns /var/log/acme-dns

# 로그
StandardOutput=journal
StandardError=journal
SyslogIdentifier=acme-dns

[Install]
WantedBy=multi-user.target
