# config/systemd/mindgarden-dev.service
[Unit]
Description=MindGarden Development Server
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/mindgarden-dev
ExecStart=/opt/mindgarden/start.sh
Restart=always
RestartSec=10
TimeoutStartSec=120
SyslogIdentifier=mindgarden-dev
# stdout/stderr 모두 journal + 파일에 남기기 (journalctl과 dev-error.log에서 예외 확인 가능)
StandardOutput=append:/var/log/mindgarden/dev.log
StandardError=append:/var/log/mindgarden/dev-error.log

[Install]
WantedBy=multi-user.target


