# scripts/automation/monitoring/log_watcher.service
[Unit]
Description=MindGarden BE Log Watcher (Discord alert)
Documentation=https://github.com/beta0629/MindGarden/blob/develop/docs/%EC%9A%B4%EC%98%81%EB%B0%98%EC%98%81/LOG_WATCHER_DISCORD_ALERT_SETUP_GUIDE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# DISCORD_WEBHOOK_URL 등 비밀값은 EnvironmentFile 로 주입한다 (코드 하드코딩 금지).
# /etc/mindgarden/log_watcher.env 파일을 chmod 600 + chown mindgarden:mindgarden 으로 보호.
EnvironmentFile=-/etc/mindgarden/log_watcher.env
# journalctl -u <unit> 권한이 있는 사용자로 실행. mindgarden 사용자가 systemd-journal 그룹 멤버여야 함.
# 권한 불충분 시 sudoers 에 NOPASSWD 로 journalctl 만 허용하는 패턴도 가능.
User=mindgarden
Group=mindgarden
ExecStart=/opt/mindgarden/scripts/automation/monitoring/log_watcher.sh

# 실패해도 timer 가 다음 주기에 재시도 — Restart=on-failure 는 oneshot 과 호환 안 됨
SuccessExitStatus=0
TimeoutStartSec=120

# 보안 강화 (선택)
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mindgarden /tmp

[Install]
WantedBy=multi-user.target
